# system/pg_backup.timer
#
# RG-Soft
#
# Тайм юнит бэкапирования баз для регламента обслуживания баз кластера PostgreSQL
#

[Unit]
Description=PG. Backup daily maintenance schedule service for databases
Requires=pg_backup.service

[Timer]
Unit=pg_backup.service
OnCalendar=Mon..Sat *-*-* 02:00:00

[Install]
WantedBy=timers.target
