# system/pg_repack.timer
#
# RG-Soft
#
<<<<<<< HEAD
# Тайм юнит очистки и перепаковки для регламента обслуживания баз кластера PG
=======
# Тайм юнит очистки и перепаковки для регламента обслуживания баз кластера PostgreSQL
>>>>>>> origin/main
#

[Unit]
Description=PG. Repack maintenance schedule service for databases
Requires=pg_repack.service

[Timer]
Unit=pg_repack.service
OnCalendar=Mon..Sat *-*-* 00:01:00

[Install]
WantedBy=timers.target
