# system/pg_vacuum_daily.service
#
# RG-Soft
#
# Сервис юнит для регламента обслуживания баз кластера PostgreSQL
#

[Unit]
Description=PG. Vacuum daily maintenance schedule service for databases

[Service]
Type=oneshot
User=postgres
#ExecStart=/var/lib/postgresql/pg_scripts/vacuum_cluster_5432_daily.sh
ExecStart=/var/lib/postgresql/pg_scripts/vacuum_dbase_zakupki_daily.sh
