# system/pg_vacuum_daily.timer
#
# RG-Soft
#
<<<<<<< HEAD
# Тайм юнит для регламента vacuum ежедневного обслуживания баз кластера PG
=======
# Тайм юнит для регламента vacuum ежедневного обслуживания баз кластера PostgreSQL
>>>>>>> origin/main
#

[Unit]
Description=PG. Vacuum daily maintenance schedule service for databases
Requires=pg_vacuum_daily.service

[Timer]
Unit=pg_vacuum_daily.service
OnCalendar=Mon..Fri *-*-* 02:00:00

[Install]
WantedBy=timers.target
